# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/dnsmasq.list
/etc/dnsmasq.conf
/etc/config/dhcp
/etc/hotplug.d/iface/13-dnsmasq
/etc/init.d/dnsmasq
/usr/sbin/dnsmasq
