# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/gw6c.list
/etc/init.d/gw6c
/etc/config/gw6c
/usr/sbin/gw6c_radvd
/usr/share/gw6c/template/openwrt.sh
/usr/sbin/gw6c
