# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/igmpproxy.list
/etc/init.d/igmpproxy
/etc/hotplug.d/iface/25-igmpproxy
/usr/sbin/igmpproxy
/sbin/act_igmp
/etc/config/igmpproxy
/etc/igmpproxy.conf
