# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/intfGrp_v6.list
/etc/init.d/default_lan_radvd
/etc/config/default_lan_radvd
