# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/iptables.list
/usr/lib/iptables/libip6t_hbh.so
/lib/firewall/iptables_app_order
/etc/config/parental
/etc/config/portTrigger
/usr/lib/iptables/libip6t_mh.so
/sbin/parental_monitor
/usr/lib/iptables/libipt_webstr.so
/etc/config/nat
/sbin/sendto.sh
/etc/hotplug.d/iface/11-nat
/sbin/firewall
/sbin/nat
/etc/init.d/firewall
/etc/config/firewall6
/lib/firewall/configure_v6filter_rules
/etc/config/route
/etc/init.d/parental_monitor
/sbin/firewall6
/usr/lib/iptables/libip6t_ipv6header.so
/lib/firewall/firewall_guest_wlan
/etc/init.d/portTrigger
/etc/init.d/nat
/usr/lib/iptables/libip6t_icmp6.so
/etc/init.d/webstr_filter
/usr/lib/firewall6.awk
/etc/firewall.user
/usr/sbin/iptables-multi
/lib/firewall/configure_firewall6
/usr/lib/iptables/libip6t_frag.so
/etc/hotplug.d/iface/13-iptables
/usr/lib/iptables/libip6t_ah.so
/etc/init.d/firewall6
/etc/init.d/route
/usr/lib/iptables/libxt_connlimit.so
/sbin/resume.sh
/sbin/webstr_filter
/usr/lib/iptables/libipt_PORTTRIGGER.so
/usr/sbin/iptables
/lib/firewall/configure_firewall
/usr/lib/iptables/libip6t_rt.so
/lib/firewall/ip6tables_app_order
/usr/lib/iptables/libip6t_eui64.so
/sbin/portTrigger
/lib/firewall/configure_filter_rules
/usr/lib/firewall.awk
