# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/lldpd.list
/etc/config/lldpd
/usr/lib/liblldpctl.so.0.0.0
/usr/sbin/lldpd
/usr/lib/liblldpctl.so
/etc/init.d/lldpd
/usr/lib/liblldpctl.so.0
/usr/sbin/lldpctl
