# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/luci-admin-core.list
/usr/lib/lua/luci/view/top.htm
/usr/lib/lua/luci/view/about.htm
/usr/lib/lua/luci/view/leftmaintenance.htm
/etc/init.d/luci_fixtime
/usr/lib/lua/luci/view/header.htm
/usr/lib/lua/luci/view/footer.htm
/usr/lib/lua/luci/view/error404.htm
/lib/uci/schema/default/dropbear
/lib/uci/schema/default/luci_hosts
/lib/uci/schema/default/wireless
/lib/uci/schema/default/system
/usr/lib/lua/luci/view/leftconfig.htm
/usr/lib/lua/luci/tools/webadmin.lua
/lib/uci/schema/default/httpd
/usr/lib/lua/luci/view/sysauth.htm
/usr/lib/lua/luci/view/path.htm
/lib/uci/schema/default/network
/usr/lib/lua/luci/view/leftmonitor.htm
/etc/init.d/luci_dhcp_migrate
/lib/uci/schema/default/fstab
/sbin/luci-flash
/usr/lib/lua/luci/view/under.htm
/usr/lib/lua/luci/view/error500.htm
/usr/lib/lua/luci/view/redirect_index.htm
/lib/uci/schema/default/dhcp
/usr/lib/lua/luci/view/indexer.htm
