# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/miniupnpd.list
/usr/sbin/miniupnpd
/sbin/act_miniupnpd
/etc/config/upnpd
/etc/hotplug.d/iface/14-miniupnpd
/etc/init.d/miniupnpd
