# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/ntpclient.list
/etc/config/ntpclient
/etc/hotplug.d/iface/20-ntpclient
/usr/sbin/adjtimex
/sbin/time_daemon
/usr/sbin/rate.awk
/usr/sbin/ntpclient
/etc/init.d/time_daemon
