# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/procd.list
/sbin/procd
/lib/functions/procd.sh
/etc/init.d/procd
