# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/proftpd.list
/etc/config/proftpd
/etc/proftpd.conf
/usr/sbin/in.proftpd
/usr/bin/ftpcount
/usr/bin/ftpwho
/usr/bin/ftptop
/etc/init.d/proftpd
/usr/sbin/proftpd
/usr/sbin/ftpshut
