# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/RA_status.list
/bin/RA_status
/etc/config/RA_status
/etc/init.d/RA_status
