# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/samba3.list
/etc/samba/smbpasswd
/etc/samba/smb.conf.template
/etc/init.d/samba
/lib/upgrade/keep.d/samba3
/usr/sbin/smbpasswd
/usr/sbin/smbd
