# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/tr143.list
/etc/config/tr143
/etc/init.d/tr143
/bin/HttpDiag
