# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/updatedd.list
/etc/config/updatedd
/etc/hotplug.d/iface/30-updatedd
/etc/init.d/updatedd
/usr/bin/updatedd
