# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/wide-dhcpv6-server.list
/usr/sbin/dhcp6s
/etc/hotplug.d/iface/50-dhcp6s
/etc/init.d/dhcp6s
/etc/config/dhcp6s
/etc/dhcp6s.conf.sample
