# Desktop/EMG2926/_V1.00(AAVK.4)C0.bin.extracted/jffs2-root/fs_1/usr/lib/opkg/info/wol.list
/etc/config/wol
/etc/hotplug.d/iface/19-wol
/sbin/startWol
/usr/bin/wolipq
/etc/init.d/wol
/usr/bin/wol
